# ly/etc/systemd/system/ly@.service.d/setfont.conf
# Ly runs on a fresh TTY (TTYReset=yes). vconsole.conf alone is unreliable here.
# Force console font before ly-dm starts. Change FONT if you prefer another size.
[Service]
ExecStartPre=-/usr/bin/setfont default8x9
